# src/bootstrap-box/home-jlee/all_configs.var.spool.cron.postgrs.cron
* * * * * /usr/local/www/cmd/util/update_cache.php >>/var/log/cmdlog_cache 2>&1
0 6 * * * /usr/local/php/bin/php /usr/local/www/cmd/util/cron.php >>/var/log/cmdlog 2>&1
0 0 * * * /usr/local/www/cmd/util/backup.sh >>/var/log/cmdlog_backups 2>&1
